# data/fwupdtool-udev/fwupdtool.service.in
[Unit]
Description=Updates devices using fwupdtool
Requires=sysinit.target
After=sysinit.target

[Service]
ExecStart=@libexecdir@/fwupd/fwupdtool update
